# src/gnugold/etc/xinetd.d/gnugol
# default: off
# description: An xinetd service that runs the gnugol protocol The
# This is the tcp version which doesn't exist yet
service gnugol
{
	disable		= yes
	id		= gnugol
        server                  = /usr/sbin/gnugold
	socket_type	= stream
	protocol	= tcp
	user		= root
	wait		= no
        flags                   = IPv6
}                                                                               

# Once invoked we want this to stay around for a while
# does mdns even work?

service gnugol
{
	id 			= gnugol
        socket_type             = dgram
        protocol                = udp
        wait                    = yes
        user                    = root
        server                  = /usr/sbin/gnugold
        disable                 = no
        per_source              = 11
        cps                     = 100 2
        flags                   = IPv6
	mdns			= yes
}
